# docker/atualizar_supervisor.conf
[program:atualizar_supervisor]
command=/opt/docker/etc/supervisor.d/atualizar_supervisor.sh
autostart=true
autorestart=true