# python/double_ia/blaze_doublerobo.conf
[program:blaze_sinais]
command=php /var/www/jonbet_usr/data/www/blaze.doublerobo.com.br/cmd.php "class=TDoubleSinaisPublisher&method=run&plataforma=blaze&idioma=ptBR"
autostart=true
autorestart=true
numprocs=1
stderr_logfile=/var/www/jonbet_usr/data/www/blaze.doublerobo.com.br/logs/blaze_sinais.err.log
stdout_logfile=/var/www/jonbet_usr/data/www/blaze.doublerobo.com.br/logs/blaze_sinais.out.log
group=blaze

[program:blaze_blaze_canal]
command=php /var/www/jonbet_usr/data/www/blaze.doublerobo.com.br/cmd.php "class=TDoubleCanalConsumer&method=run&canal_id=1"
autostart=true
autorestart=true
numprocs=1
stderr_logfile=/var/www/jonbet_usr/data/www/blaze.doublerobo.com.br/logs/blaze_blaze_canal.err.log
stdout_logfile=/var/www/jonbet_usr/data/www/blaze.doublerobo.com.br/logs/blaze_blaze_canal.out.log
group=blaze

[program:blaze_telegram]
command=php /var/www/jonbet_usr/data/www/blaze.doublerobo.com.br/cmd.php "class=TTelegramConsumer&method=run"
autostart=true
autorestart=true
numprocs=5
process_name=%(program_name)s_%(process_num)02d
stderr_logfile=/var/www/jonbet_usr/data/www/blaze.doublerobo.com.br/logs/blaze_telegram_%(process_num)s.err.log
stdout_logfile=/var/www/jonbet_usr/data/www/blaze.doublerobo.com.br/logs/blaze_telegram_%(process_num)s.out.log
group=blaze

[program:blaze_historico]
command=php /var/www/jonbet_usr/data/www/blaze.doublerobo.com.br/cmd.php "class=TDoubleHistoricoConsumer&method=run"
autostart=true
autorestart=true
numprocs=1
stderr_logfile=/var/www/jonbet_usr/data/www/blaze.doublerobo.com.br/logs/blaze_historico.err.log
stdout_logfile=/var/www/jonbet_usr/data/www/blaze.doublerobo.com.br/logs/blaze_historico.out.log
group=blaze

[program:blaze_historico_usuario]
command=php /var/www/jonbet_usr/data/www/blaze.doublerobo.com.br/cmd.php "class=TDoubleUsuarioHistoricoConsumer&method=run"
autostart=true
autorestart=true
numprocs=5
process_name=%(program_name)s_%(process_num)02d
stderr_logfile=/var/www/jonbet_usr/data/www/blaze.doublerobo.com.br/logs/blaze_historico_usuario_%(process_num)s.err.log
stdout_logfile=/var/www/jonbet_usr/data/www/blaze.doublerobo.com.br/logs/blaze_historico_usuario_%(process_num)s.out.log
group=blaze

[program:blaze_cron_atualiza_objetivos]
command=php /var/www/jonbet_usr/data/www/blaze.doublerobo.com.br/cmd.php "class=TDoubleCron&method=atualizar_objetivos"
autostart=true
autorestart=true
numprocs=1
stderr_logfile=/var/www/jonbet_usr/data/www/blaze.doublerobo.com.br/logs/blaze_cron_atualiza_objetivos.err.log
stdout_logfile=/var/www/jonbet_usr/data/www/blaze.doublerobo.com.br/logs/blaze_cron_atualiza_objetivos.out.log
group=blaze

[program:blaze_notification_server]
command=php /var/www/jonbet_usr/data/www/blaze.doublerobo.com.br/app/service/rest/websocket/NotificationServer.php
autostart=true
autorestart=true
numprocs=1
stderr_logfile=/var/www/jonbet_usr/data/www/blaze.doublerobo.com.br/logs/blaze_notification_server.err.log
stdout_logfile=/var/www/jonbet_usr/data/www/blaze.doublerobo.com.br/logs/blaze_notification_server.out.log
group=blaze

[program:blaze_redis_websocket]
command=php /var/www/jonbet_usr/data/www/blaze.doublerobo.com.br/cmd.php "class=RedisWebSocket&method=run&plataforma=blaze&idioma=ptBR"
autostart=true
autorestart=true
numprocs=1
stderr_logfile=/var/www/jonbet_usr/data/www/blaze.doublerobo.com.br/logs/blaze_redis_websocket.err.log
stdout_logfile=/var/www/jonbet_usr/data/www/blaze.doublerobo.com.br/logs/blaze_redis_websocket.out.log
group=blaze
